# config/cyberdeck.service
[Unit]
Description=CyberDeck AP + MQTT boot setup
After=network.target NetworkManager.service
Wants=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cyberdeck-start-ap.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
