# config/cyberdeck-status.service
[Unit]
Description=CyberDeck status collector
After=cyberdeck.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/collect-status.sh
